# deploy/clawstack.service.example
[Unit]
Description=ClawStack backend
After=network.target

[Service]
Type=simple
User=clawstack
Group=clawstack
WorkingDirectory=/opt/clawstack
Environment=CLAWSTACK_ENV_FILE=/opt/clawstack/backend/.env
ExecStart=/opt/clawstack/.venv/bin/python -m uvicorn clawstack_backend.app:app --app-dir /opt/clawstack/backend --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
